# rocm-persistenced.service
[Service]
ExecStart=/opt/rocm-persistenced/bin/rocm-persistenced

[Install]
WantedBy=multi-user.target
